# orbita-bot.service
[Unit]
Description=ORBITA AI Telegram Bot (Clean Architecture with DI)
After=network.target postgresql.service

[Service]
Type=simple
User=alex
WorkingDirectory=/home/alex/apps/tg-analytics-bot
Environment="PATH=/home/alex/apps/tg-analytics-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/alex/apps/tg-analytics-bot"
Environment="USE_DI_HANDLERS=true"
ExecStart=/home/alex/apps/tg-analytics-bot/venv/bin/python -m app.main_di
Restart=always
RestartSec=10
StandardOutput=append:/home/alex/apps/tg-analytics-bot/logs/bot.log
StandardError=append:/home/alex/apps/tg-analytics-bot/logs/bot-error.log

[Install]
WantedBy=multi-user.target
